# jailer/flavours/sample/usr/local/etc/cron.d/pkglist
0 0 * * * 		root pkg prime-origins > /var/backups/pkg.list
